# systemd/docker-image-tracker.service
[Unit]
Description=Docker Image Usage Tracker
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /path/to/runner-cleanup/track-docker-image-usage.py --state /var/lib/runner-cleanup/image-usage.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
